# cookbooks/duckduckhack-vm/files/default/etc/xdg/autostart/update-duckpan.desktop
[Desktop Entry]
Name=DuckPAN Update
Type=Application
Terminal=true
Exec=/usr/local/bin/update-duckpan.sh
